# Source/cron.d/root
*/1 * * * * sh /data/adb/modules/BATTERYOPT/WHITELIST.sh
